# csv_syncer/csv_syncer.crontab
*/5 * * * * /usr/bin/python3 /root/csv_syncer/csv_syncer.py >> /root/csv_syncer/listener.log 2>&1
